# crontab/reboot-system.sh
#!/bin/bash
#
# Script para ejecutar aplicacion Node reboot-system
#

# 1. Carga el PATH donde se encuentra reboot-system
cd /path/to/node/reboot-system
# 2. Ejecuta la aplicacion de Node reboot-system
node app.js $1
